# systemd/trustpanel-expire.service
[Unit]
Description=TrustPanel Expire Users + Sync
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/opt/trustpanel
EnvironmentFile=/opt/trustpanel/.env

ExecStart=/opt/trustpanel/venv/bin/python -m services.expire_service

User=root
Group=root