# systemd/trustsystem-webhook.service
[Unit]
Description=TrustSystem Webhook
After=network.target

[Service]
User=root
WorkingDirectory=/opt/trustsystem
ExecStart=/opt/trustsystem/webhook_venv/bin/python run_webhook.py
Restart=always

[Install]
WantedBy=multi-user.target