# systemd/trustsystem-admin.service
[Unit]
Description=TrustSystem Admin Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/trustsystem
EnvironmentFile=/opt/trustsystem/.env
ExecStart=/opt/trustsystem/venv/bin/python -m bots.admin.bot

Restart=always
RestartSec=5
User=root

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
