# systemd/trustsystem-public.service
[Unit]
Description=TrustSystem Public Bot
After=network.target

[Service]
Type=simple
WorkingDirectory=/opt/trustsystem

ExecStart=/bin/bash -c 'set -a; source /opt/trustsystem/.env; set +a; exec /opt/trustsystem/venv/bin/python -m bots.public.bot'

Restart=always
RestartSec=5
User=root

StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
